# rc.kfcgi
#!/bin/sh

# PROVIDE: kfcgi
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name=kfcgi
rcvar=kfcgi_enable
load_rc_config ${name}

: ${kfcgi_enable:=NO}
kfcgi_flags="${kfcgi_flags} -- ${kfcgi_prog}"

command=/usr/local/sbin/kfcgi

run_rc_command "$1"
